# bin/kmotion@.service
[Unit]
Description=kmotion
Requires=network.target
After=network.target

[Service]
ExecStart=/usr/local/kmotion3/kmotion.py
User=%I
Restart=always
RestartSec=60

[Install]
WantedBy=multi-user.target
