# bin/kmotion.service
[Unit]
Description=kmotion
Requires=network.target
After=network.target

[Service]
ExecStart=/usr/local/kmotion3/kmotion.py
User=videouser

[Install]
WantedBy=multi-user.target
