# bin/reboot_cams.timer
[Timer]
OnCalendar=01:17

[Install]
WantedBy=timers.target