# bin/reboot_cams.service
[Unit]
Description=kmotion reboot cameras

[Service]
Type=oneshot
User=videouser
ExecStart=/usr/local/kmotion3/reboot_cams.py reboot
